# meta-bsp/recipes-kb/kb-regulatory/kb-regulatory/kb-regulatory.service
[Unit]
Description=User Space Regulatory Firmware Loading

[Service]
ExecStart=/etc/kb-regulatory

[Install]
WantedBy=multi-user.target
